# kitchen_manager_addon/rootfs/etc/s6-overlay/s6-rc.d/kitchen_manager/run
#!/bin/bash

# Simple bash run script - avoids s6-overlay-suexec PID 1 check
echo "Starting Kitchen Manager (Simple Bash)..."

export DATA_PATH="/data"

cd /app
exec python3 -m uvicorn src.main:app --host 0.0.0.0 --port 8000
